# supervisor/supervisord_headphones.conf
[program:headphones]
directory=/opt/headphones
command=python Headphones.py
user = root
autostart=true
autorestart=true